# station-checkout/deploy/passport-checkout.service
[Unit]
Description=Passport Checkout Station
After=network.target

[Service]
Type=simple
User=passport
Group=passport
WorkingDirectory=/opt/passport-system/station-checkout
Environment="PATH=/opt/passport-system/station-checkout/.venv/bin"
Environment="PYTHONPATH=/opt/passport-system"
ExecStart=/opt/passport-system/station-checkout/.venv/bin/python -m app.daemon
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# NFC + GPIO přístup
SupplementaryGroups=plugdev dialout i2c gpio

[Install]
WantedBy=multi-user.target
